# routing/myrouting.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/myrouting.sh

[Install]
WantedBy=default.target
